# _superninja_startup.conf

[program:8050_python]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec python -m http.server 8050"
directory=/workspace/mechanic-shop
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/8050_python.err.log
stdout_logfile=/var/log/supervisor/8050_python.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2


[program:9000_python]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec python -m http.server 9000"
directory=/workspace/mechanic-shop
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/9000_python.err.log
stdout_logfile=/var/log/supervisor/9000_python.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2


[program:8090_python3]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec python3 -m http.server 8090"
directory=/workspace/mechanic-shop
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/8090_python3.err.log
stdout_logfile=/var/log/supervisor/8090_python3.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2


[program:9084_github-mcp-serv]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec ./github-mcp-server"
directory=/opt/ninja-mcp-servers/github-mcp
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/9084_github-mcp-serv.err.log
stdout_logfile=/var/log/supervisor/9084_github-mcp-serv.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2


[program:18080_psiphon-tunnel-]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec /opt/ninja-browser-use/psiphon-tunnel-core-x86_64 -config /opt/ninja-browser-use/psiphon.config.json"
directory=/opt/ninja-browser-use
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/18080_psiphon-tunnel-.err.log
stdout_logfile=/var/log/supervisor/18080_psiphon-tunnel-.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2


[program:8085_python3]
command=/bin/bash -c "source ~/.bashrc 2>/dev/null; exec python3 -m http.server 8085"
directory=/workspace/mechanic-shop
user=root
autostart=true
autorestart=true
startsecs=3
startretries=3
stderr_logfile=/var/log/supervisor/8085_python3.err.log
stdout_logfile=/var/log/supervisor/8085_python3.out.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

